# dot_dotfiles/setup/roles/pi/avahi/assets/etc/systemd/system/avahi-aliases.service
[Unit]
Description=avahi subdomain announcer
After=network.target

[Service]
ExecStart=/usr/local/bin/avahi-alias start
ExecStopPost=/usr/local/bin/avahi-alias stop

[Install]
WantedBy=multi-user.target
